# com.ibc.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ibc.frontend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/serve</string>
        <string>-s</string>
        <string>dist</string>
        <string>-l</string>
        <string>3000</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/william/factory/inspired-benevolence-cybersecurity/frontend</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/william/factory/inspired-benevolence-cybersecurity/frontend/logs/frontend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/william/factory/inspired-benevolence-cybersecurity/frontend/logs/frontend-error.log</string>
</dict>
</plist>
