# com.ibc.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ibc.backend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>dist/index.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/william/factory/inspired-benevolence-cybersecurity/backend</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/william/factory/inspired-benevolence-cybersecurity/backend/logs/backend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/william/factory/inspired-benevolence-cybersecurity/backend/logs/backend-error.log</string>
</dict>
</plist>
